# systemd/user/wallpaper-update.service
[Unit]
Description=Update wallpaper

[Service]
Type=oneshot
ExecStart=/home/mugiwara/.config/scripts/reloads/reload_wallpaper.sh
